# https-outcall.service
[Unit]
Description=HTTPS Outcall Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/https_outcall
ExecStart=/root/https_outcall/target/release/https-outcall
Restart=on-failure
RestartSec=5
StartLimitBurst=3
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target 