# https-outcall.service
[Unit]
Description=HTTPS Outcall Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/https_outcall
ExecStart=/root/https_outcall/target/release/https-outcall
Restart=on-failure
RestartSec=5
StartLimitBurst=3
Environment=RUST_LOG=info
# SSL is now handled by DigitalOcean load balancer
# Environment=USE_HTTPS=true
# Environment=SSL_CERT_PATH=/root/https_outcall/certs/cert.pem
# Environment=SSL_KEY_PATH=/root/https_outcall/certs/key.pem

[Install]
WantedBy=multi-user.target 